# script.timer
[Unit]
Description=Run script dayly and on boot

[Timer]
OnBootSec=1min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
